# upstart/docker-mk-svn.conf
description "Subversion container"
author "michael Kueper <michael@kueper-datentechnik.de>
start on filesystem and started docker
stop on runlevel [!2345]
respawn
script
    /usr/bin/docker start mk-svn
end script
